# devops/configs/supervisord.ini
[program:gunicorn]
command=/usr/local/bin/gunicorn app:app -c /vagrant/gunicorn.conf.py
directory=/vagrant
user=root
autostart=true
autorestart=true
redirect_stderr=true
